# rclone-systemd/gmedia.mount
[Unit]
Description = /gmedia MergerFS mount
After=gmedia-rclone.service
RequiresMountsFor=/data

[Mount]
What = /data/local:/GD
Where = /gmedia
Type = fuse.mergerfs
Options = sync_read,use_ino,allow_other,auto_cache,splice_move,splice_read,splice_write,func.getattr=newest,category.action=all,category.create=ff

[Install]
WantedBy=multi-user.target
